# deploy/systemd/robin-digest-daily.timer
[Unit]
Description=Robin: daily digest at 09:00 Asia/Tbilisi (duties.md TZ)

[Timer]
OnCalendar=*-*-* 09:00:00 Asia/Tbilisi
Persistent=true

[Install]
WantedBy=timers.target
