# deploy/systemd/robin-liveness.timer
[Unit]
Description=Robin: hourly liveness check

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
