# deploy/systemd/robin-digest-daily.service
[Unit]
Description=Robin: compose and post the daily ecosystem digest
After=robin-mirror-sync.service

[Service]
Type=oneshot
User=robin
WorkingDirectory=/srv/robin/robin-runtime
EnvironmentFile=/srv/robin/robin.env
ExecStart=/srv/robin/robin-runtime/.venv/bin/python -m robin.digest daily
